# etc/udev/rules.d/00-backlight.rules
SUBSYSTEM=="backlight",  KERNEL=="intel_backlight", ACTION=="add", RUN+="/bin/chgrp users %S%p/brightness", RUN+="/bin/chmod 664 %S%p/brightness"
